# launchd/com.notifyme.checker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.notifyme.checker</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Update this path to your Python installation -->
        <string>/usr/local/bin/python3</string>
        <string>-m</string>
        <string>notifyme.cli</string>
        <string>check</string>
    </array>

    <!-- Check every 5 minutes (300 seconds) -->
    <key>StartInterval</key>
    <integer>300</integer>

    <!-- Start immediately when loaded -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/Users/mooney/.notifyme/scheduler.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mooney/.notifyme/scheduler-error.log</string>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/mooney/NotifyMe</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/mooney/NotifyMe</string>
    </dict>
</dict>
</plist>
